# ATM.Machine/hardware/system/aaeon/atm_base-machine.conf
exec /opt/launchers/start-atm_base-machine
respawn
start on (local-filesystems and net-device-up IFACE!=lo)
chdir /opt/apps/machine/atm_base-machine
